# configs/gl-mt300nv2-dlnasrv/rootfs_overlay/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
sleep 1
ntfs-3g /dev/sda1 /mnt/sda1 -o rw,lazytime,noatime,big_writes
exit 0
